# crontab/generateEdm.command
#!/bin/sh
PATH=/usr/local/bin:/usr/local/sbin:~/bin:/usr/bin:/bin:/usr/sbin:/sbin
DIR = 'your E-newsletter folder path'
#ex: ~/Work/Newsletter/EDMgenerator

cd DIR && yarn run edm